# home/.config/autostart/flameshot.desktop
[Desktop Entry]
Name=Flameshot
Comment=A tool to take screenshots.
Type=Application
Exec=flameshot
StartupNotify=false
Terminal=false
Hidden=false
