# home/.config/autostart/picom.desktop
[Desktop Entry]
Name=Picom
Comment=
Type=Application
Exec=picom -b
StartupNotify=false
Terminal=false
Hidden=false
