# home/.config/autostart/plank.desktop
[Desktop Entry]
Name=Plank Dock
Comment=A stupidly simple dock.
Type=Application
Exec=plank
StartupNotify=false
Terminal=false
Hidden=false
