# home/.config/autostart/redshift.desktop
[Desktop Entry]
Name=Redshift
Comment=Blue light filter.
Type=Application
Exec=redshift
StartupNotify=false
Terminal=false
Hidden=false
